# supervisor/ftpserver.ini
[program:mainftpserver]
command=python ftpserver/ftpserver.py
environment=FLASK_ENV="production"
process_name=mainftpserver
numprocs=1
directory=/var/www/hosting/code

startretries=3
autorestart=true
exitcodes=0

stdout_logfile=/var/www/hosting/log/ftpserver.out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

stderr_logfile=/var/www/hosting/log/ftpserver.err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=10
