# supervisor/hosting.ini
uwsgi --plugin=python --http=127.0.0.1:8005 --wsgi-file=src/__init__.py --need-app
[program:$supervisor_program_name]
command=uwsgi /var/www/hosting/supervisor/hosting.uwsgi.ini

process_name=main:hosting
numprocs=1
directory=/var/www/hosting/code

startretries=3
autorestart=true
exitcodes=0

stdout_logfile=/var/www/hosting/log/out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

stderr_logfile=/var/www/hosting/log/err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=10
