# supervisor/programtemplate.ini
[program:$supervisor_program_name]
command=firejail --private-bin=bash,fish,python,pip python $entrypoint

process_name=$supervisor_program_name
numprocs=1
directory=$program_ftpdir

startretries=3
autorestart=true
exitcodes=0

stdout_logfile=./programs/$username/$program_name.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=10

stderr_logfile=./programs/$username/$program_name.err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=10
